# siglab/systemd/siglab-report-push.timer
[Unit]
Description=Push siglab's signal report every 15 min (report writer isn't clock-aligned, so this just re-pushes whatever's newest on disk each firing)

[Timer]
OnCalendar=*-*-* *:00/15:00
Persistent=true

[Install]
WantedBy=timers.target
