# siglab/systemd/siglab-report-push.timer
[Unit]
Description=Check every 15 min for a new siglab daily digest to push (2026-07-24: detailed hourly reports stay local-only, see push_report.sh)

[Timer]
OnCalendar=*-*-* *:00/15:00
Persistent=true

[Install]
WantedBy=timers.target
